# etc/files/crontab
*/55  *  *  *  *    salt '*' state.highstate
